# rootfs/etc/services.d/pgbouncer/run
#!/usr/bin/with-contenv bash

pgbouncer -u postgres /etc/pgbouncer/pgbconf.ini