# rootfs/etc/service/pgbouncer/run
#!/bin/bash

exec /usr/sbin/pgbouncer -u postgres /etc/pgbouncer/pgbconf.ini