# rootfs/etc/service/pgbouncer/run
#!/bin/bash

exec /usr/sbin/pgbouncer /etc/pgbouncer/pgbconf.ini